# packaging/systemd/diskmon.service
[Unit]
Description=Diskmon periodic run (one-shot)
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/etc/diskmon
# Use wrapper to capture logs and rotate; change path if necessary
ExecStart=/usr/local/bin/diskmon-run.sh
StandardOutput=journal
StandardError=journal
Nice=10
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
